# dot_config/xremap/systemd/xremap.service
[Unit]
Description=xremap service

[Service]
Restart=always
ExecStart=/usr/local/bin/xremap %h/.config/xremap/config.yml

[Install]
WantedBy=default.target
