# config/xremap/systemd/xremap.service
[Unit]
Description=xremap service

[Service]
Restart=always
ExecStart=/usr/bin/xremap ~/.config/xremap/config.yml

[Install]
WantedBy=default.target
